# wimz-poweroff-pulse.service
[Unit]
Description=Pulse GPIO26 to kill Pololu latch at shutdown
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target
Requires=shutdown.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wimz_poweroff_pulse.sh
TimeoutStartSec=10
RemainAfterExit=yes

[Install]
WantedBy=halt.target reboot.target shutdown.target
