# xbox-controller.service
[Unit]
Description=WIM-Z Xbox Controller Service
After=network.target bluetooth.target
Wants=network.target

[Service]
Type=simple
User=morgan
Group=morgan
WorkingDirectory=/home/morgan/dogbot
Environment=PYTHONPATH=/home/morgan/dogbot
ExecStart=/home/morgan/dogbot/env_new/bin/python /home/morgan/dogbot/xbox_hybrid_controller.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Health monitoring
WatchdogSec=60
TimeoutStopSec=15
KillMode=mixed

[Install]
WantedBy=multi-user.target