# scripts/systemd/wimz-healthcheck.timer
[Unit]
Description=Run WIM-Z liveness healthcheck every 60s

[Timer]
OnBootSec=2min
OnUnitActiveSec=60s
AccuracySec=5s

[Install]
WantedBy=timers.target
